# glyphbox.service
[Unit]
Description=GLYPHBOX Fantasy Console
# Wait for basic system initialisation — sound, input, DRM are all needed
After=sound.target local-fs.target systemd-udev-settle.service
Wants=sound.target

[Service]
Type=simple
User=greedo
Group=greedo

WorkingDirectory=/home/greedo/glyphbox

# SDL2 on Pi without a desktop environment uses the KMS/DRM driver directly.
# AUDIODEV/SDL_AUDIODRIVER ensure ALSA is used rather than PulseAudio.
Environment=SDL_VIDEODRIVER=kmsdrm
Environment=SDL_AUDIODRIVER=alsa
Environment=HOME=/home/greedo

# Give the DRM subsystem a moment to settle before SDL2 opens the display.
ExecStartPre=/bin/sleep 2

ExecStart=/home/greedo/glyphbox/build/glyphbox

# Restart automatically if glyphbox crashes (not on clean exit / intentional quit)
Restart=on-failure
RestartSec=3

# If it fails repeatedly, stop trying after 5 attempts within 60 seconds
StartLimitIntervalSec=60
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
